# infra/php/crontab
0 6,18 * * * /usr/local/bin/php /var/www/html/bin/console app:update-weather-cache && redis-cli BGSAVE >> /var/www/html/var/log/cron.log 2>&1
